# templates/default/ssh_tunnel.conf.erb
description     "ssh_tunnel to <%= @remote_host %>"

start on net-device-up
stop on deconfiguring-networking

respawn

exec sudo -u <%= node["ssh_tunnel"]["source_user"] %> ssh -N <%= @tunnels %> <%= @username %>@<%= @remote_host %>

